# packaging/linux/agent-up-server.service
[Unit]
Description=Agent-Up Server
After=network.target

[Service]
Type=simple
ExecStart=/opt/agent-up/server/AgentUp.Server
Environment=ASPNETCORE_URLS=http://localhost:5000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
