# packaging/linux/agent-up-server.service
[Unit]
Description=Agent-Up Server
After=network.target

[Service]
Type=simple
ExecStart=/opt/agent-up/server/AgentUp.Server --urls http://127.0.0.1:5000
Environment=ASPNETCORE_URLS=http://127.0.0.1:5000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
